# templates/default/mongodb.upstart.erb
# THIS FILE HAS BEEN GENERATED BY CHEF
# ANY MANUAL MODIFICATIONS WILL BE OVERWRITTEN

description "MongoDB v<%= node[:mongodb][:version] %>"

kill timeout <%= node[:mongodb][:timeout] %> # wait this many seconds between SIGTERM and SIGKILL.

start on runlevel [2345]
stop on runlevel [!2345]

respawn
respawn limit 10 5

exec sudo -u mongodb /usr/bin/mongod -f /etc/mongodb.conf
